# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/etc.conf
[Service]
ExecStartPre=cp -rpn /usr/share/factory/etc /
