# mkosi.images/base/mkosi.extra/etc/rc.local
#!/bin/sh -ue

# Fix services imported by systemd-confext in /etc not being loaded.
systemctl daemon-reload
# This is intended to fix some problems caused by systemd-sysext.
# When systemd-sysext merge sysext images, it won't trigger those services inside to be activated.
systemctl --no-block default
