# mkosi.images/default-initrd/mkosi.extra/usr/lib/systemd/system/systemd-confext-sysroot.service.d/exec.conf
[Service]
ExecStart=
ExecStart=systemd-confext --root=/sysroot refresh --noexec=no --mutable=import
